# TN_docker/files/TN_docker
#!/bin/sh /etc/rc.common

START=98

start() {
	sleep 30
	dockerd & 
}

stop() {
	killall dockerd 
}

